# stutter.service
[Unit]
Description=Focus-aware process priority daemon for Hyprland
PartOf=graphical-session.target
After=graphical-session.target

[Service]
ExecStart=%h/.local/bin/stutter
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=3s

[Install]
WantedBy=graphical-session.target